# deploy/systemd/ai-radar-enricher.service
[Unit]
Description=ai-radar enricher
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/workspace/ai-radar
EnvironmentFile=%h/workspace/ai-radar/.env
ExecStart=%h/workspace/ai-radar/enricher/build/install/enricher/bin/enricher
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
